# ec2/postgres/postgres-s3-backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Daily backup at 1:00 AM (server time — UTC on EC2)
0 1 * * * root /usr/local/sbin/pg-backup-to-s3.sh >> /var/log/postgres-backup.log 2>&1

# Daily cleanup at 2:00 AM (after backup window)
0 2 * * * root /usr/local/sbin/pg-backup-cleanup-s3.sh >> /var/log/postgres-backup.log 2>&1
